# web/supervisord.conf
[program:node]
command = node app.js 80
directory = /home/dotcloud/current/
